# deploy/qibla.service
[Unit]
Description=Travel Smart Travel Fast (Qibla) Node server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/qibla
EnvironmentFile=/var/www/qibla/.env
ExecStart=/usr/bin/node server/index.js
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/qibla

[Install]
WantedBy=multi-user.target
